# soc-mediatek-mt6795/soc-mediatek-mt6795-enable-wlan0.initd
#!/sbin/openrc-run

name="enable-wlan0"
description="Mediatek MT6795 Enable Wlan0"
command=/usr/sbin/enable-wlan0
pidfile="/run/enable-wlan0.pid"
command_background=true