# soc-mediatek-mt6795/soc-mediatek-mt6795-autokd.initd
#!/sbin/openrc-run

name="soc-mediatek-mt6795-autokd"
description="Mediatek MT6795 Autok Daemon"
command=/usr/sbin/autokd
pidfile="/run/autokd.pid"
command_background=true